# systemd/meshcore-discovery.timer
[Unit]
Description=Run MeshCore discovery daily at 2:00 AM

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
